# deploy/raspberry-pi/wall-clock.service
[Unit]
Description=Wall Clock local server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/wall-clock/wall-clock -addr :8080 -config /etc/wall-clock/config.json -public-url http://wallclock.local:8080/control
Restart=always
RestartSec=2
User=pi
Group=pi
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=read-only
ReadWritePaths=/etc/wall-clock
WorkingDirectory=/opt/wall-clock

[Install]
WantedBy=multi-user.target
